# deploy/warlock-server.service
# systemd unit for the Warlock multiplayer WebSocket server.
# Install:  sudo cp deploy/warlock-server.service /etc/systemd/system/
#           sudo systemctl daemon-reload && sudo systemctl enable --now warlock-server
# Logs:     journalctl -u warlock-server -f

[Unit]
Description=Warlock multiplayer server
After=network.target

[Service]
Type=simple
# Adjust to where you cloned the repo and which user runs it.
WorkingDirectory=/opt/warlock
User=warlock
Environment=NODE_ENV=production
Environment=PORT=8787
# tsx runs the TypeScript server directly (no build step). Uses the local copy
# from node_modules so the unit doesn't depend on a global install.
ExecStart=/usr/bin/env node_modules/.bin/tsx server/index.ts
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
